# tools/ci/server/deploy/systemd/ci-jenkins.service
[Unit]
Description=CI Jenkins (offline WAR + JCasC)
After=network.target

[Service]
Type=simple
# 密钥（JENKINS_ADMIN_PASSWORD）由 deploy.py 从 config.local.ini 写入，0600，不入仓（C-1）。
EnvironmentFile=@@ENV_FILE@@
Environment=JENKINS_HOME=@@JENKINS_HOME@@
Environment=CASC_JENKINS_CONFIG=@@CASC@@
ExecStart=@@JAVA@@ -Djenkins.install.runSetupWizard=false -jar @@WAR@@ --httpPort=@@HTTP_PORT@@
Restart=always
RestartSec=5
# Jenkins 首次解插件启动较慢，放宽启动超时。
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
